# deploy/nanoclaw.service
[Unit]
Description=NanoClaw Agent
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/root/nanoclaw
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
EnvironmentFile=/root/nanoclaw/.env
Environment=HOME=/root
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
