# .claude/skills/add-watchdog/setup/watchdog.timer
[Unit]
Description=Run NanoClaw watchdog every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Unit=nanoclaw-watchdog.service

[Install]
WantedBy=timers.target
