# porcupine/com.nanoclaw.voice-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.nanoclaw.voice-daemon</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/pavelkudrna/nanoclaw/porcupine/venv/bin/python3</string>
		<string>/Users/pavelkudrna/nanoclaw/porcupine/voice-daemon.py</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/pavelkudrna/nanoclaw</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONUNBUFFERED</key>
		<string>1</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/pavelkudrna/nanoclaw/logs/voice-daemon.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/pavelkudrna/nanoclaw/logs/voice-daemon.error.log</string>
</dict>
</plist>
