# scripts/sqlite-web.service
[Unit]
Description=sqlite-web — read-only browser for NanoClaw DB
After=network.target nanoclaw.service
BindsTo=nanoclaw.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/nanoclaw
EnvironmentFile=/home/ubuntu/nanoclaw/.env
ExecStart=/bin/sh -c 'sqlite_web --host 0.0.0.0 --port 8088 --read-only --password "${SQLITE_WEB_PASSWORD}" store/messages.db'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
