# scripts/podcast/com.nanoclaw.local-enrichment-nightly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nanoclaw.local-enrichment-nightly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/ilansolot/nanoclaw-v2/scripts/podcast/run-local-enrichment-nightly.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/ilansolot/nanoclaw-v2</string>

  <key>StandardOutPath</key>
  <string>/Users/ilansolot/nanoclaw-v2/logs/local-enrichment-nightly.launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ilansolot/nanoclaw-v2/logs/local-enrichment-nightly.launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
