# launchd/com.nanoclaw.token-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nanoclaw.token-refresh</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/gf/nanoclaw/scripts/refresh-token.mjs</string>
    </array>
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>StandardOutPath</key>
    <string>/Users/gf/nanoclaw/logs/token-refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gf/nanoclaw/logs/token-refresh.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
